# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

ln -nsf /mnt/sda1/portal/wifidog /www/wifidog
exit 0
